# net-misc/agocontrol/files/agosimulator.init
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

command="/opt/agocontrol/bin/agosimulator.py"
pidfile="/var/run/${SVCNAME}.pid"
start_stop_daemon_args="--background --make-pidfile --chdir $(dirname ${command})"
